# chezmoi-state/private_dot_config/systemd/user/syncthing.service
[Unit]
Description=Syncthing
Requires=home.mount

[Service]
ExecStart=/usr/bin/syncthing

[Install]
WantedBy=default.target
